# scripts/lab-control.service
[Unit]
Description=Lab Control System
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/opt/lab-control
ExecStart=/usr/bin/python3 /opt/lab-control/launcher.py --daemon
ExecStop=/usr/bin/python3 /opt/lab-control/launcher.py --stop
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lab-control

[Install]
WantedBy=multi-user.target
